# systemd/tova-ui.service
[Unit]
Description=TOVA UI Gunicorn daemon
After=network.target postgresql.service
Requires=postgresql.service

[Service]
# User and group to run as
User=www-data
Group=www-data

# Working directory
WorkingDirectory=/Users/danielstephens/Desktop/TOVA/ui

# Environment variables
Environment="PATH=/Users/danielstephens/Desktop/TOVA/ui/venv/bin"
Environment="FLASK_SECRET_KEY=your-secret-key-here"
Environment="DATABASE_URL=postgresql+psycopg2://user:password@localhost:5432/mydb"
Environment="POSTGRES_HOST=localhost"
Environment="POSTGRES_PORT=5432"
Environment="POSTGRES_USER=user"
Environment="POSTGRES_PASSWORD=password"
Environment="POSTGRES_DB=mydb"
Environment="GUNICORN_BIND=0.0.0.0:8000"
Environment="GUNICORN_WORKERS=4"
Environment="GUNICORN_LOG_LEVEL=info"
Environment="RUN_DB_INIT_ON_START=1"

# ExecStart command
ExecStart=/Users/danielstephens/Desktop/TOVA/ui/venv/bin/gunicorn \
    --config /Users/danielstephens/Desktop/TOVA/ui/gunicorn_config.py \
    app:server

# Restart policy
Restart=always
RestartSec=10

# Security settings
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/Users/danielstephens/Desktop/TOVA

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tova-ui

[Install]
WantedBy=multi-user.target

